# aethercloud-v6.service
[Unit]
Description=AetherCloud DynamicV6 SOCKS5 gateway
Wants=network-online.target docker.service
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/aethercloud-v6 run
ExecStop=/usr/local/sbin/aethercloud-v6 stop
Restart=always
RestartSec=5s
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
